# roles/scrutiny/templates/scrutiny.timer.j2
[Unit]
Description=Scrutiny disk check

[Timer]
OnCalendar={{ scrutiny_schedule }}
Persistent=true

[Install]
WantedBy=timers.target
