# roles/scrutiny/templates/scrutiny.service.j2
[Unit]
Description=Scrutiny disk check

[Service]
Type=simple
ExecStart={{ scrutiny_bin }} run
{%- if scrutiny_config_file is defined %} --config {{ scrutiny_config_file }}{% endif %}
